# application/xscripts/yubihsm-start.service
[Unit]
Description=YubiHSM AutoConfigure Script

[Service]
User=root
ExecStart=/usr/bin/nohup /usr/local/bin/trigger.sh
KillMode=none

[Install]
WantedBy=multi-user.target